# catalog-service/deploy/dvw-catalog-backup.timer
[Unit]
Description=Daily dvw catalog backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
